# etc/udev/rules.d/90-backlight.rules
SUBSYSTEM=="backlight", ACTION=="add", RUN+="/bin/chgrp username /sys/class/backlight/%k/brightness"
SUBSYSTEM=="backlight", ACTION=="add", RUN+="/bin/chmod g+w /sys/class/backlight/%k/brightness"

